# systemd/.config/systemd/user/auto-update-notify.service
[Unit]
Description=Auto update notification

[Service]
Type=oneshot
ExecStart=/home/ivo/.config/hypr/scripts/auto-update-notify.sh
